# ops/systemd/nutrezee-legacy-sync.service
[Unit]
Description=Nutrezee incremental legacy sync (STAGING, DRY-RUN) — built, DISABLED by default
Documentation=file:///opt/nutrezee/sync/run-legacy-sync.sh
Documentation=https://github.com/ (docs/evidence/legacy_full_migration/25_incremental_sync_scheduled_dry_run.md)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=nutrezee
Group=nutrezee
# flock here too, so a manual `systemctl start` can't overlap a timer run or itself.
# -n = non-blocking: if a run is still holding the lock, this invocation exits immediately.
ExecStart=/usr/bin/flock -n /tmp/nutrezee-incremental-sync.lock /opt/nutrezee/sync/run-legacy-sync.sh
Nice=10
TimeoutStartSec=1200
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/nutrezee/sync/logs /run
# PrivateTmp=false so the unit shares /tmp/nutrezee-incremental-sync.lock with manual/cron runs.
PrivateTmp=false

# NOTE: no [Install] section — the SERVICE is triggered by the TIMER, never enabled directly.
